# huawei-monitoring.timer
[Unit]
Description=Startup connection status monitoring

[Timer]
OnBootSec=1min
OnUnitActiveSec=60s

[Install]
WantedBy=timers.target