# huawei-monitoring.service
[Unit]
Description=Connection status monitoring

[Service]
ExecStart=/usr/bin/huawei-umnc -m
